# net_monitor.upstart
author "Alexander Grothe <seahawk1986@gmx.de>"
description "Avahi linker"

start on started vdr
stop on stopping vdr

respawn

script
exec /usr/bin/net_monitor
end script

