# update_monitor.upstart
author "Alexander Grothe <seahawk1986@gmx.de>"
description "Avahi linker"

start on started avahi-linker or started vdr
stop on stopping avahi-linker or stopping vdr

respawn

script
exec /usr/bin/update_monitor -w /srv/vdr/video.00/.update
end script

